# 0x1A-application_server/gunicorn.service
[Unit]
Description=gunicorn daemon

[Service]
Type=simple
ExecStart=/usr/bin/env bash /etc/init.d/gunicorn_sd.sh
